# install/wacli-sync-personal.service.template
[Unit]
Description=wacli WhatsApp Sync (Personal)
After=network.target

[Service]
Type=simple
ExecStart={{BINARY_PATH}} --store {{STORE_PATH}} sync
Restart=always
RestartSec=10
Environment="HOME={{HOME}}"
StandardOutput=append:{{HOME}}/.wacli/sync.log
StandardError=append:{{HOME}}/.wacli/sync.err

[Install]
WantedBy=default.target
