# install/wacli-sync.service.template
[Unit]
Description=wacli WhatsApp Sync (__INSTANCE_NAME__)
After=network.target

[Service]
Type=simple
ExecStart=__BINARY_PATH__ --store __STORE_PATH__ sync
Restart=always
RestartSec=10
Environment="HOME=__HOME__"
StandardOutput=append:__STORE_PATH__/sync.log
StandardError=append:__STORE_PATH__/sync.err

[Install]
WantedBy=default.target
